# supervisor-progs/trytond.conf
[program:trytond]
command=/usr/local/bin/trytond -c /etc/trytond.conf
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true
